# packaging/deb/postinst
#!/bin/sh
# Demeter .deb post-install: create the service user, seed an admin on first
# install, then enable + (re)start the systemd service.
set -e

SERVICE=demeter
RUN_USER=demeter
DATA_DIR=/var/lib/demeter
BIN=/usr/bin/demeter
CRED="$DATA_DIR/INITIAL_ADMIN_PASSWORD"

case "$1" in
	configure)
		# Dedicated system user.
		if ! getent passwd "$RUN_USER" >/dev/null; then
			adduser --system --group --no-create-home --home "$DATA_DIR" \
				--shell /usr/sbin/nologin "$RUN_USER" >/dev/null
		fi
		mkdir -p "$DATA_DIR"
		chown "$RUN_USER:$RUN_USER" "$DATA_DIR"
		chmod 0750 "$DATA_DIR"

		# First install (empty data dir): create an admin login. If this fails for
		# any reason, the service's own first-run bootstrap still generates one and
		# writes the same INITIAL_ADMIN_PASSWORD file, so the password is never lost.
		PASS=""
		if [ -z "$(ls -A "$DATA_DIR" 2>/dev/null)" ]; then
			PASS="$(LC_ALL=C tr -dc 'A-Za-z0-9' </dev/urandom 2>/dev/null | head -c 16 || true)"
			if runuser -u "$RUN_USER" -- "$BIN" --create-admin "admin:${PASS}" --data-dir "$DATA_DIR" >/dev/null 2>&1; then
				{ echo "username: admin"; echo "password: $PASS"; } > "$CRED" 2>/dev/null || true
				chown "$RUN_USER:$RUN_USER" "$CRED" 2>/dev/null || true
				chmod 0600 "$CRED" 2>/dev/null || true
			else
				PASS="" # fall back to the service's bootstrap (writes $CRED on start)
			fi
		fi

		systemctl daemon-reload >/dev/null 2>&1 || true
		systemctl enable "$SERVICE" >/dev/null 2>&1 || true
		systemctl restart "$SERVICE" >/dev/null 2>&1 || true

		IP="$(hostname -I 2>/dev/null | awk '{print $1}')"
		echo ""
		echo "=================== Demeter installed ==================="
		echo " URL:    http://${IP:-<server-ip>}:8080"
		if [ -n "$PASS" ]; then
			echo " Login:  admin / $PASS"
		fi
		echo " First-login password is saved here (root only):"
		echo "   sudo cat $CRED"
		echo " Change it in the GUI, then delete that file."
		echo " Status: systemctl status $SERVICE   Logs: journalctl -u $SERVICE"
		echo "========================================================"
		;;
esac

exit 0
